# supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=debug               ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[program:flaskapp1]
command=/home/vagrant/Envs/flaskapp1/bin/python flaskapp1/flask-app1.py ;
directory=/vagrant/ ;
stderr_logfile=/vagrant/logs/flaskapp1.stderr ;
environment=PORT=8000;

[program:flaskapp2]
command=/home/vagrant/Envs/flaskapp2/bin/python flaskapp2/flask-app2.py ;
directory=/vagrant/ ;
stderr_logfile=/vagrant/logs/flaskapp2.stderr ;
environment=PORT=8001;

[inet_http_server]
port=0.0.0.0:9001   ;

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://0.0.0.0:9001 ;
